# pkg/debian/postinst
#!/bin/bash
ldconfig
systemctl enable loxilb
systemctl start loxilb
